# security/softethervpn/files/vpnclient.in
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: vpnclient
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add these lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# vpnclient_enable (bool):	Set to NO by default.
#				Set it to YES to enable vpnclient.

. /etc/rc.subr

name=vpnclient
rcvar=vpnclient_enable
start_cmd="${name}_start"

load_rc_config $name

: ${vpnclient_enable:="NO"}

command=%%DATADIR%%/${name}/${name}

vpnclient_start() {
	${command} start
}

run_rc_command "$1"
