# ftp/proftpd-devel/pkg-plist
bin/ftpasswd
bin/ftpcount
bin/ftpdctl
bin/ftpmail
bin/ftpquota
bin/ftptop
bin/ftpwho
bin/prxs
%%ETCDIR%%.conf.sample
etc/rc.d/proftpd
include/proftpd/ascii.h
include/proftpd/auth.h
include/proftpd/bindings.h
include/proftpd/buildstamp.h
include/proftpd/child.h
include/proftpd/class.h
include/proftpd/cmd.h
include/proftpd/compat.h
include/proftpd/conf.h
include/proftpd/config.h
include/proftpd/configdb.h
include/proftpd/ctrls.h
include/proftpd/data.h
include/proftpd/default_paths.h
include/proftpd/dirtree.h
include/proftpd/display.h
include/proftpd/encode.h
include/proftpd/env.h
include/proftpd/event.h
include/proftpd/expr.h
include/proftpd/feat.h
include/proftpd/filter.h
include/proftpd/fsio.h
include/proftpd/ftp.h
include/proftpd/glibc-glob.h
include/proftpd/help.h
include/proftpd/ident.h
include/proftpd/inet.h
include/proftpd/lastlog.h
include/proftpd/libsupp.h
include/proftpd/log.h
include/proftpd/memcache.h
include/proftpd/mkhome.h
include/proftpd/mod_ctrls.h
include/proftpd/mod_log.h
include/proftpd/mod_tls.h
include/proftpd/modules.h
include/proftpd/netacl.h
include/proftpd/netaddr.h
include/proftpd/netio.h
include/proftpd/options.h
include/proftpd/parser.h
include/proftpd/pidfile.h
include/proftpd/pool.h
include/proftpd/pr-syslog.h
include/proftpd/privs.h
include/proftpd/proctitle.h
include/proftpd/proftpd.h
include/proftpd/regexp.h
include/proftpd/response.h
include/proftpd/rlimit.h
include/proftpd/scoreboard.h
include/proftpd/session.h
include/proftpd/sets.h
include/proftpd/signals.h
include/proftpd/stash.h
include/proftpd/str.h
include/proftpd/support.h
include/proftpd/table.h
include/proftpd/throttle.h
include/proftpd/timers.h
include/proftpd/tpl.h
include/proftpd/trace.h
include/proftpd/utf8.h
include/proftpd/var.h
include/proftpd/version.h
include/proftpd/xferlog.h
libdata/pkgconfig/proftpd.pc
man/man1/ftpasswd.1.gz
man/man1/ftpcount.1.gz
man/man1/ftpmail.1.gz
man/man1/ftpquota.1.gz
man/man1/ftptop.1.gz
man/man1/ftpwho.1.gz
man/man5/proftpd.conf.5.gz
man/man5/xferlog.5.gz
man/man8/ftpdctl.8.gz
man/man8/ftpscrub.8.gz
man/man8/ftpshut.8.gz
man/man8/proftpd.8.gz
sbin/ftpscrub
sbin/ftpshut
sbin/in.proftpd
sbin/proftpd
@dir lib/proftpd
@dir %%LOCALSTATEDIR%%
