# sysutils/puppet/files/puppetmaster.in
#!/bin/sh
#
# $FreeBSD$
#

# PROVIDE: puppetmaster
# REQUIRE: NETWORK

# Add the following lines to /etc/rc.conf to enable the puppetmaster:
#
# puppetmaster_enable="YES"
#

. /etc/rc.subr

name="puppetmaster"
rcvar=puppetmaster_enable

load_rc_config "${name}"

: ${puppetmaster_enable="NO"}

command="%%PREFIX%%/bin/puppet"
command_args="master"
command_interpreter=%%RUBY%%

pidfile="$($command config print --section master pidfile)"

run_rc_command "$1"
